# cmd/ical-notifier/misc/ical-notifier.cron
0,15,30,45 * * * * ical-notifier --config /etc/ical-relay/config.yml > /proc/1/fd/1 2>/proc/1/fd/2
