# cmd/ical-notifier/misc/ical-notifier.cron
*/15 * * * * ical-notifier --config /etc/ical-relay/noitifer-config.yml > /proc/1/fd/1 2>/proc/1/fd/2
